# systemd/user/delete-downloads.service
[Unit]
Description=Delete unwanted '~/Downloads' directory.

[Service]
Type=oneshot
ExecStart=/usr/bin/rm -rf %h/Downloads

