# systemd/user/delete-downloads.path
[Unit]
Description=Watch for unwanted '~/Downloads' directory

[Path]
PathExists=%h/Downloads

[Install]
WantedBy=default.target